# deploy_data/profile.d/docker-env.sh
#!/bin/sh

#*For Debian copy to /etc/profiles.d/ 


export GITLAB_HOME=/srv/gitlab
export CNT_GITLAB="gitlab-web-1"
export SESOP_HOME=/srv/sesop
